# init/hosts-update.service
[Unit]
Description=Update MVPS hosts

[Service]
Type=oneshot
ExecStart=/usr/bin/hosts-update
After=network.target
